# container/rescue.service
[Unit]
Description=Rescue Shell

[Service]
Type=oneshot
ExecStart=/bin/bash
StandardInput=tty
TTYVHangup=yes
TTYReset=yes
