# container/wyrefi.service
[Unit]
Description=Wyrcan EFI
ConditionKernelCommandLine=|wyrcan.efi=clear
ConditionKernelCommandLine=|wyr.efi=clear
Before=wyrnet.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wyrcan efi
StandardInput=tty
TTYVHangup=yes
TTYReset=yes

[Install]
RequiredBy=wyrnet.service
